# configs/examples/systemd/cs2-auto-update.service
[Unit]
Description=CS2 dedicated server auto-update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/cs2-server-ops/apps/maintain/updater/update_cs2.sh
