# configs/examples/systemd/cs2-auto-update.timer
[Unit]
Description=Run CS2 auto-update daily

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
